# ics_backend/deploy/.service
[Unit]
Description=MTUCI Rasp ICS backend
After=network.target
StartLimitIntervalSec=0

[Service]
Environment="%ENV%"
Type=simple
Restart=always
RestartSec=1
User=%USER%
ExecStart=/usr/bin/java -jar %JAR%

[Install]
WantedBy=multi-user.target